# bs/files/etc/cron.daily/certbot-renew
#!/bin/bash

certbot renew --noninteractive --renew-hook /bs/rpc-explorer/bs/bin/restart-nginx.sh
