# bs/files/etc/cron.d/update
# update script
* * * * *               root       /bs/elements-explorer/bs/bin/cibuild-cron >/tmp/update.out 2>&1
